# com.zeromd.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zeromd.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__ZEROMD_SYNC_SCRIPT__</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__ZEROMD_HOME__/.zeromd/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__ZEROMD_HOME__/.zeromd/launchd-stderr.log</string>
</dict>
</plist>
